# systemd/gamepi-kiosk.service
[Unit]
Description=GamePi13 MIDIJuggler Kiosk on tty1
Wants=gamepi-kiosk-ready.service midijuggler.service
Requires=midijuggler.service
After=gamepi-kiosk-ready.service gamepi-splash.service network-online.target midijuggler.service
Conflicts=getty@tty1.service
ConditionPathExists=!/run/gamepi-console-boot

[Service]
Type=simple
User=dietpi
SupplementaryGroups=video tty input
TimeoutStartSec=60
TimeoutStopSec=15
KillMode=control-group
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
Environment=HOME=/home/dietpi
Environment=FRAMEBUFFER=/dev/fb0
Environment=GAMEPI_SPLASH_COMPLETED_FLAG=/run/gamepi/splash-completed
RuntimeDirectory=gamepi
WorkingDirectory=/home/dietpi
TTYPath=/dev/tty1
TTYReset=no
TTYVHangup=no
StandardInput=tty
StandardOutput=journal
StandardError=journal
ExecStart=/opt/midijuggler/app/scripts/gamepi-launch-kiosk.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
