# systemd/gamepi-kiosk-ready.service
[Unit]
Description=Wait for MIDIJuggler web UI and release GamePi splash
After=network-online.target midijuggler.service gamepi-splash.service
Wants=network-online.target midijuggler.service
Before=gamepi-kiosk.service
ConditionPathExists=!/run/gamepi-console-boot

[Service]
Type=oneshot
RemainAfterExit=yes
User=dietpi
Environment=MIDIJUGGLER_SCRIPTS_DIR=/opt/midijuggler/app/scripts
Environment=GAMEPI_SPLASH_COMPLETED_FLAG=/run/gamepi/splash-completed
RuntimeDirectory=gamepi
StandardOutput=journal
StandardError=journal
ExecStartPre=+/opt/midijuggler/app/scripts/gamepi-disable-blanking.sh
ExecStart=/opt/midijuggler/app/scripts/wait-for-midijuggler-web.sh
ExecStart=+/opt/midijuggler/app/scripts/gamepi-splash-stop.sh

[Install]
WantedBy=multi-user.target
