# systemd/wing-gadget-loop.service
[Unit]
Description=Route g_audio USB gadget input to Wing USB channels 5-6
After=sound.target
Requires=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/alsaloop -C plughw:CARD=g_audio,DEV=0 -P wing_stereo3 -f S16_LE -t 5000 -b
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
