# systemd/wing-gadget-loop.service
[Unit]
Description=Route g_audio USB gadget input to Wing USB channels 5-6
After=sound.target
Wants=sound.target

[Service]
Type=simple
Environment=MIDIJUGGLER_SCRIPTS_DIR=/opt/midijuggler/app/scripts
ExecStart=/opt/midijuggler/app/scripts/wing-gadget-loop.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
